# scripts/com.chessckers.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Layer-0 supervisor for the chessckers training watchdog. The watchdog
  itself supervises the tmux sessions (workers / sync_leena / sync_vast)
  inside the user-level tmux server. This plist keeps the watchdog alive
  even if the tmux server crashes — KeepAlive=true relaunches on any exit.

  Install:
    cp scripts/com.chessckers.watchdog.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.chessckers.watchdog.plist

  Stop / uninstall:
    launchctl unload -w ~/Library/LaunchAgents/com.chessckers.watchdog.plist

  Tail log:
    tail -f /tmp/watchdog.launchd.log

  Notes:
   - The watchdog script itself exits cleanly when $LOCAL_RUN/STOP exists,
     so end-of-run shutdown is respected (launchd then won't relaunch
     because RunAtLoad fires only on load, not on every clean exit; but
     KeepAlive=true would. To pause supervision when a run ends, either
     unload the plist or use SuccessfulExit=false in KeepAlive.)
   - Edit RUN_DIR / POLL_SEC env via the dict below if needed.
   - This file path embeds an absolute repo location — adjust if the
     repo moves.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chessckers.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>/bin/bash</string>
        <string>/Users/ox/AAworkspace/chessckers/scripts/watchdog.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ox/AAworkspace/chessckers</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LOCAL_RUN</key>
        <string>/Users/ox/AAworkspace/chessckers/engine/runs/local-006</string>
        <key>POLL_SEC</key>
        <string>60</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Relaunch on any exit (clean or crash). When a run ends cleanly via
         STOP, the watchdog exits 0; launchd will relaunch it, see the STOP
         file again, and exit immediately — small busy-loop but harmless.
         To suppress that, unload the plist when the run ends. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Don't relaunch faster than every 30s if it keeps dying. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/watchdog.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/watchdog.launchd.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
